# configs/wehttamsnaps/docs/files/~/.config/systemd/user/regenerate-cheatsheet.service.d/env.conf
[Service]
# Set REPO_ROOT to your repo location using %h (the user's home directory).
# Change the path if your repository is in a different location.
Environment=REPO_ROOT=%h/WehttamSnaps-Niri