# configs/systemd/user/config-backup.service
[Unit]
Description=Daily backup of WehttamSnaps SwayFx configurations
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.config/wehttamsnaps/scripts/backup-configs.sh

[Install]
WantedBy=default.target
